# etc/systemd/system/omarchy-dynamic-function-row-daemon.service
[Unit]
Description=omarchy-dynamic-function-row-daemon - Legacy alias for tiny-dfr-omarchy
After=systemd-user-sessions.service systemd-logind.service
Conflicts=tiny-dfr-omarchy.service tiny-dfr.service

[Service]
Type=simple
User=root

# Wait for device readiness before daemon startup
ExecStartPre=-/usr/bin/udevadm settle -t 4
ExecStartPre=-/usr/bin/wait-for-device.sh touchbar-minimal 4
ExecStartPre=-/usr/bin/wait-for-device.sh hid-kbd 4

ExecStart=/usr/bin/tiny-dfr-omarchy
Restart=always
RestartSec=1

# Security settings for hardware daemon
ProtectHome=false
ProtectSystem=false
ProtectKernelModules=true
ProtectControlGroups=true
PrivateUsers=false
PrivateNetwork=true
PrivateTmp=false
RestrictSUIDSGID=false
SupplementaryGroups=users
CapabilityBoundingSet=CAP_SETUID CAP_SETGID CAP_DAC_OVERRIDE CAP_FOWNER
Environment="XDG_RUNTIME_DIR=/run/user"

[Install]
WantedBy=multi-user.target
