# etc/systemd/system/suspend-fix-t2.service
[Unit]
Description=Disable and Re-Enable Apple BCE Module + WiFi + KB Backlight
Before=sleep.target
StopWhenUnneeded=yes

[Service]
User=root
Type=oneshot
RemainAfterExit=yes

# Save keyboard backlight and turn it off during suspend.
ExecStart=-/usr/bin/tiny-dfr-kbd-backlight save-off

#touchbar and dependencies
ExecStart=-/usr/bin/systemctl stop omarchy-dynamic-function-row-daemon.service
ExecStart=-/usr/bin/modprobe -r appletbdrm
ExecStart=-/usr/bin/modprobe -r hid_appletb_bl
ExecStart=-/usr/bin/modprobe -r hid_appletb_kbd

# WiFi reload is intentionally disabled for faster wake.
# Re-enable these lines only if your WiFi fails after resume.
# ExecStart=-/usr/bin/modprobe -r brcmfmac_wcc
# ExecStart=-/usr/bin/modprobe -r brcmfmac

ExecStart=-/usr/bin/rmmod -f apple-bce

ExecStop=-/usr/bin/modprobe apple-bce

# WiFi reload is intentionally disabled for faster wake.
# Re-enable these lines only if your WiFi fails after resume.
# ExecStop=-/usr/bin/modprobe brcmfmac
# ExecStop=-/usr/bin/modprobe brcmfmac_wcc

#touchbar and dependencies - BCE-aware resume sequence
ExecStop=-/usr/bin/wait-for-device.sh touchbar-bce 10
ExecStop=-/usr/bin/modprobe hid_appletb_bl
ExecStop=-/usr/bin/wait-for-device.sh hid-backlight 4
ExecStop=-/usr/bin/modprobe hid_appletb_kbd
ExecStop=-/usr/bin/wait-for-device.sh hid-kbd 4
ExecStop=-/usr/bin/tiny-dfr-kbd-backlight restore
ExecStop=-/usr/bin/modprobe appletbdrm
ExecStop=-/usr/bin/wait-for-device.sh touchbar-minimal 4
ExecStopPost=-/usr/bin/systemctl reset-failed omarchy-dynamic-function-row-daemon.service
ExecStopPost=-/usr/bin/systemctl restart omarchy-dynamic-function-row-daemon.service

[Install]
WantedBy=sleep.target
