# etc/systemd/system/suspend-fix-t2.service
[Unit]
Description=suspend-fix-t2 - Apple BCE suspend/resume coordinator
Before=sleep.target
StopWhenUnneeded=yes

[Service]
Type=oneshot
User=root
RemainAfterExit=yes

# Suspend sequence: save backlight, stop daemons, unload Touch Bar modules, power down BCE
ExecStart=-/usr/bin/tiny-dfr-kbd-backlight save-off
ExecStart=-/usr/bin/systemctl stop tiny-dfr-omarchy.service
ExecStart=-/usr/bin/systemctl stop omarchy-dynamic-function-row-daemon.service
ExecStart=-/usr/bin/modprobe -r appletbdrm
ExecStart=-/usr/bin/modprobe -r hid_appletb_bl
ExecStart=-/usr/bin/modprobe -r hid_appletb_kbd
ExecStart=-/usr/bin/rmmod -f apple-bce

# Resume sequence: reinitialize BCE, reload Touch Bar modules, restore services
# (WiFi reload disabled for faster wake - uncomment brcmfmac if needed)
ExecStop=-/usr/bin/modprobe apple-bce
ExecStop=-/usr/bin/wait-for-device.sh touchbar-bce 3
ExecStop=-/usr/bin/modprobe hid_appletb_bl
ExecStop=-/usr/bin/wait-for-device.sh hid-backlight 4
ExecStop=-/usr/bin/modprobe hid_appletb_kbd
ExecStop=-/usr/bin/wait-for-device.sh hid-kbd 4
ExecStop=-/usr/bin/tiny-dfr-kbd-backlight restore
ExecStop=-/usr/bin/modprobe appletbdrm
ExecStop=-/usr/bin/wait-for-device.sh touchbar-minimal 4

# Let udev finish device registration before daemon rebinds
ExecStopPost=-/usr/bin/udevadm settle -t 3
ExecStopPost=-/usr/bin/systemctl reset-failed tiny-dfr-omarchy.service
ExecStopPost=-/usr/bin/systemctl restart tiny-dfr-omarchy.service

[Install]
WantedBy=sleep.target
